# systemd/octave-snapshot-daily.timer
[Unit]
Description=Daily octave-snapshot build of GNU Octave

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
